# deploy/supervisor.conf
; ═══════════════════════════════════════════════════════════════════════════════
; A3 ERP — Configuration Supervisor (queue workers)
; ═══════════════════════════════════════════════════════════════════════════════
; Installation :
;   sudo apt install supervisor
;   sudo cp deploy/supervisor.conf /etc/supervisor/conf.d/a3erp.conf
;   sudo supervisorctl reread && sudo supervisorctl update
;   sudo supervisorctl start a3erp-worker:*
;
; Vérification :
;   sudo supervisorctl status
; ═══════════════════════════════════════════════════════════════════════════════

[program:a3erp-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home/votre-domaine.com/public_html/artisan queue:work redis --sleep=3 --tries=3 --max-time=3600 --queue=default,emails,notifications,pdf
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/log/supervisor/a3erp-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=120

; Worker dédié aux PDFs lourds (timeout plus long)
[program:a3erp-pdf-worker]
process_name=%(program_name)s_%(process_num)02d
command=php /home/votre-domaine.com/public_html/artisan queue:work redis --sleep=5 --tries=2 --max-time=7200 --queue=pdf --timeout=300
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/log/supervisor/a3erp-pdf-worker.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=3
stopwaitsecs=300

[group:a3erp]
programs=a3erp-worker,a3erp-pdf-worker
